# systemd/telinha.service
[Unit]
Description=Project250 Telinha ESP32
After=multi-user.target

[Service]
Type=simple
User=project250
WorkingDirectory=/opt/project250
ExecStart=/usr/bin/python3 /opt/project250/telinha.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
